# support/author-service/agent-loop-author.socket
[Unit]
Description=Agent Loop fixed Codex author broker socket

[Socket]
ListenStream=/run/agent-loop/author.sock
Accept=yes
SocketMode=0600
SocketUser=root
SocketGroup=root
DirectoryMode=0755
RemoveOnStop=yes
MaxConnections=4

[Install]
WantedBy=sockets.target
