# support/author-service/agent-loop-author@.service
[Unit]
Description=Agent Loop fixed Codex author broker connection
Requires=agent-loop-author.socket
After=agent-loop-author.socket
CollectMode=inactive-or-failed

[Service]
Type=exec
ExecStart=/usr/bin/python3 -I -B -S -c "import importlib.machinery as m,runpy,sys;root='/opt/agent-loop-author-service/runtime/lib/python3.14/site-packages';package=root+'/agent_loop';source_loader=(m.SourceFileLoader,m.SOURCE_SUFFIXES);sys.path_importer_cache[root]=m.FileFinder(root,source_loader);sys.path_importer_cache[package]=m.FileFinder(package,source_loader);sys.path.insert(0,root);runpy.run_module('agent_loop.author_service',run_name='__main__')"
StandardInput=socket
StandardOutput=socket
StandardError=journal
Environment=AGENT_LOOP_AUTHOR_BROKER_UNIT=%n
EnvironmentFile=/etc/agent-loop/author-service.conf
User=root
Group=root
UMask=0077
# Host-visible, root-private publication point for the immutable closure copies
# that PID 1 binds into the unprivileged author unit.  Individual broker
# requests remove their own UUID-named directory after proving cgroup emptiness;
# preserving the empty parent avoids races between concurrent socket instances.
RuntimeDirectory=agent-loop/author-closures
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
TemporaryFileSystem=/runtime:nodev,nosuid,size=536870912,mode=0700
NoNewPrivileges=yes
# The broker must traverse the already-open, operator-owned mode-0700
# CODEX_HOME descriptor to witness it.  It receives no path authority and no
# write capability; every other capability remains outside the bounding set.
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_DAC_READ_SEARCH
PrivateDevices=yes
PrivateTmp=disconnected
PrivateIPC=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=all
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native
KillMode=control-group
SendSIGKILL=yes
TimeoutStopSec=5s
OOMPolicy=kill
RuntimeMaxSec=46min
MemoryMax=1536M
TasksMax=64
LimitCORE=0
LimitFSIZE=384M
LimitNOFILE=256
